# build/rq_worker/Worker/supervisord.conf
[supervisord]
nodaemon=true

[program:worker]
command=python3 worker.py
process_name=%(program_name)s-%(process_num)s
numprocs=4
directory=.
stopsignal=TERM
autostart=true
autorestart=true
